# rootfs/etc/s6-overlay/s6-rc.d/sendspin/run
#!/command/with-contenv bash
# S6 longrun service: start the Sendspin Bridge Python application.
exec /app/entrypoint.sh
